# edge-door-server/edge-door-server.timer
[Unit]
Description=Edge Door Server

[Timer]
OnCalendar=*:*:0/20

[Install]
WantedBy=timers.target
